# bin/crontab
41 3 * * * certbot renew --noninteractive --post-hook "systemctl restart mbta"