# bin/mbta.service
[Unit]
Wants=network-online.target
After=network-online.target

[Service]
WorkingDirectory=/home/ubuntu
ExecStart=/usr/bin/bash /home/ubuntu/mbta_service/mbta/bin/run_service.sh
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=MBTA
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target
